# examples/debian/pm_ingest.service
[Unit]
Description=Send PatchMe version data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pm_ingest.sh
